# demo-password-reset.service
[Unit]
Description=Reset Django demo user password
After=network.target

[Service]
Type=oneshot
User=storage
WorkingDirectory=/home/storage/Cloude/cloudservice
ExecStart=/bin/bash /home/storage/Cloude/reset_demo_password.sh
